# crontab
0 12 * * *  /usr/local/bin/python /usr/src/scraper-task/script/db_dump.py >> /tmp/cron.log 2>&1